# telegram-privilege-monitor.service
[Unit]
Description=Telegram Privilege Monitor - Surveillance des élévations de privilèges
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/telegram_notif/telegram_privilege_monitor.sh start
ExecStop=/usr/local/bin/telegram_notif/telegram_privilege_monitor.sh stop
Restart=always
RestartSec=10
User=root
Group=root

# Sécurité
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/lib/telegram_privilege_monitor /var/run /tmp
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram_privilege_monitor

[Install]
WantedBy=multi-user.target
