# telegram-privilege-monitor.service
[Unit]
Description=Telegram Privilege Monitor - Surveillance des élévations de privilèges
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/telegram_notif/telegram_privilege_monitor.sh start
ExecStop=/usr/local/bin/telegram_notif/telegram_privilege_monitor.sh stop
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sécurité
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/lib/telegram_privilege_monitor /var/run

# Variables d'environnement
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CURL_TIMEOUT=10

[Install]
WantedBy=multi-user.target
